# postgresql-patroni/postgresql-patroni.d/local/share/cook/templates/patroni.rc.in
#!/bin/sh
#
# PROVIDE: patroni
# REQUIRE: NETWORKING LOGIN FILESYSTEM

. /etc/rc.subr

name="patroni"
rcvar=patroni_enable

list_cmd="patroni_list"
history_cmd="patroni_history"
extra_commands="list history reload"

load_rc_config patroni
# load postgres configuration for user and group permission
load_rc_config postgresql
: ${patroni_enable:=no}
: ${patroni_config:="/usr/local/etc/patroni/patroni.yml"}
: ${patroni_pid:="/tmp/patroni.pid"}
: ${patroni_cmd:="/usr/local/bin/patroni"}
: ${patroni_ctl:="/usr/local/bin/patronictl"}
: ${postgresql_user:="postgres"}

pidfile=${patroni_pid}

procname="$(sed -n '1s/^#!//p' "$patroni_cmd")"
command="/usr/sbin/daemon"
command_args="-f -c -S -u ${postgresql_user} -p ${patroni_pid} \
    -t ${name} -T ${name} ${patroni_cmd} ${patroni_config}"

patroni_list()
{
	${patroni_ctl} -c ${patroni_config} list
}

patroni_history()
{
	${patroni_ctl} -c ${patroni_config} history
}

run_rc_command "$1"
