# zincsearch/zincsearch.d/local/share/cook/templates/zincsearch.rc.in
#!/bin/sh

# PROVIDE: zincsearch
# REQUIRE: LOGIN
# KEYWORD: shutdown

# Add the following lines to /etc/rc.conf to enable zincsearch
# zincsearch_enable="YES"
#
# zincsearch_enable (bool):
#     Set it to YES to enable zincsearch
#     Set to NO by default
# zincsearch_user (string):
#     Set user that zincsearch will run under
#     Default is "zincsearch"
# zincsearch_group (string):
#     Set group that own zincsearch files
#     Default is "zincsearch"
# zincsearch_syslog_output_enable (bool)
#     Set it to YES to enable stdout/stderr syslog output
#     Set to NO by default
# zincsearch_syslog_output_tag (str)
#     Set syslog tag if syslog enabled
#     Default is "zincsearch"
# zincsearch_args (string)
#     Set additional command line arguments
#     Default is ""

. /etc/rc.subr

name=zincsearch
rcvar=zincsearch_enable

load_rc_config $name

: ${zincsearch_enable:="NO"}
: ${zincsearch_user:="zincsearch"}
: ${zincsearch_group:="zincsearch"}
: ${zincsearch_syslog_output_enable:="NO"}
: ${zincsearch_args:=""}

if checkyesno zincsearch_syslog_output_enable; then
	if [ -n "${zincsearch_syslog_output_tag}" ]; then
		zincsearch_syslog_output_flags="-T ${zincsearch_syslog_output_tag}"
	else
		zincsearch_syslog_output_flags="-T ${name}"
	fi
fi

pidfile="/var/run/${name}/${name}.pid"
required_files=""

procname="/usr/local/bin/zincsearch"
command="/usr/sbin/daemon"
command_args="-f ${zincsearch_syslog_output_flags} -p ${pidfile} -t ${name} \
 /usr/bin/env ${zincsearch_env} ZINC_FIRST_ADMIN_USER=%%zincuser%% \
 ZINC_FIRST_ADMIN_PASSWORD=%%zincpass%% ZINC_DATA_PATH=%%zincdata%% \
 ZINC_SERVER_PORT=%%zincport%% ZINC_PROMETHEUS_ENABLE=True \
 ZINC_PLUGIN_ES_VERSION=8.8 \
 ${procname} \
 ${zincsearch_args}"

start_precmd="zincsearch_start_precmd"

zincsearch_start_precmd()
{
	if [ ! -d "/var/run/${name}" ]; then
		install -d -m 0750 -o ${zincsearch_user} -g ${zincsearch_group} "/var/run/${name}"
	fi
}

run_rc_command "$1"